# kiosk_skeleton/etc/systemd/system/kiosk-schedule-screen.service
# /etc/systemd/system/kiosk-schedule-screen.service
[Unit]
Description=Schedule Screen On/Off Timers
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/schedule-screen-timers

[Install]
WantedBy=multi-user.target
